# able/able.service
[Unit]
Description=ABLE - Autonomous Business & Learning Engine
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=able
Group=able
WorkingDirectory=/opt/able/ABLE
Environment="PATH=/home/able/.able/venv/bin:/home/able/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="ABLE_HOME=/home/able/.able"

# Load environment from .env file
EnvironmentFile=-/home/able/.able/.env

# Run the ABLE agent
ExecStart=/home/able/.able/venv/bin/able

# Restart policy — back off on crashes, don't burn CPU
Restart=on-failure
RestartSec=30

# Memory guard — prevent OOM on 2GB VPS (leave ~512MB for OS)
MemoryMax=1536M
MemoryHigh=1280M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=able

[Install]
WantedBy=multi-user.target
